# deploy/systemd/cloudflared-watchdog.service.example
[Unit]
Description=Cloudflared connectivity watchdog check
After=cloudflared.service
Wants=cloudflared.service

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '/usr/bin/curl --silent --show-error --fail --max-time 8 https://pplinsighthub.com/health > /dev/null || /usr/bin/systemctl restart cloudflared.service'
