# deploy/systemd/cloudflared-watchdog.timer.example
[Unit]
Description=Run cloudflared watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Unit=cloudflared-watchdog.service

[Install]
WantedBy=timers.target
